# src/rootfs_skel/etc/inittab
# Initial system critical filesystems mounting (/proc, /dev, /sys, etc.)
::sysinit:/sbin/openrc sysinit

# Goto other runlevels now
::wait:/sbin/openrc boot
::wait:/sbin/openrc default
::wait:/sbin/openrc

# Put a getty on the serial port
ttyS0::respawn:/sbin/getty -L ttyS0 115200 vt100

# Logging junk
null::sysinit:/bin/touch /var/log/messages
null::respawn:/sbin/syslogd -n -m 0
null::respawn:/sbin/klogd -n

# What to do at the "Three Finger Salute"
::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
null::shutdown:/sbin/openrc shutdown
null::shutdown:/bin/umount -a -r
null::shutdown:/sbin/swapoff -a

